# world/rootfs/etc/rc.local
#!/bin/sh

# Add your own startup tasks here
